# crontab
TZ=Asia/Tokyo
SHELL=/bin/bash

# 土曜日と水曜の朝9時
0 9 * * 3,6 lp -d printer -o fitplot /test-print.pdf > /proc/1/fd/1 2>&1
